# config/systemd/dbus.socket
[Unit]
Description=D-Bus System Message Bus Socket

[Socket]
ListenStream=/opt/app/platform/avahi/var/run/dbus/system_bus_socket
