# config/systemd/platform-nginx.service
[Unit]
Description=Platform nginx
Requires=platform-uwsgi-internal.service platform-uwsgi-public.service
After=platform-uwsgi-internal.service platform-uwsgi-public.service

[Service]
Type=forking
PIDFile=/opt/app/platform/nginx/logs/nginx.pid
#RootDirectory=/opt/app/platform/nginx
ExecStartPre=/opt/app/platform/nginx/sbin/nginx -t -c /opt/app/platform/config/nginx/nginx.conf
ExecStart=/opt/app/platform/nginx/sbin/nginx -c /opt/app/platform/config/nginx/nginx.conf
ExecReload=/opt/app/platform/nginx/sbin/nginx -c /opt/app/platform/config/nginx/nginx.conf -s reload
ExecStop=/opt/app/platform/nginx/sbin/nginx -c /opt/app/platform/config/nginx/nginx.conf -s stop

[Install]
WantedBy=multi-user.target